# dep_search-docker/crontab
SHELL=/bin/bash
@daily  cd /root && bash update.bosque.sh
